# marca.conf
description "node.js server"
author      "kvz - http://kevin.vanzonneveld.net"

# Used to Be: Start on Startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
#stop on shutdown

start on (net-device-up
          and local-filesystems
      and runlevel [2345])
stop on runlevel [016]

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    #export HOME="/home/multiservice"
    export HOME="/"
    #exec /usr/local/bin/node /opt/datos/server.js >> /var/log/datos.log 2>&1
    cd /opt/marca2
    NODE_ENV=production exec /usr/bin/node /opt/marca2/server.js >> /var/log/marca.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
